# deployment/shared/systemd/certus-neo4j.service
[Unit]
Description=Certus Neo4j Graph Database
After=network-online.target podman.service
Wants=network-online.target
Requires=podman.service

[Service]
Type=simple
User=certus
Restart=always
RestartSec=10s
TimeoutStartSec=300

EnvironmentFile=/opt/certus/.env

ExecStartPre=-/usr/bin/podman stop certus-neo4j
ExecStartPre=-/usr/bin/podman rm certus-neo4j
ExecStartPre=/usr/bin/podman pull docker.io/library/neo4j:5-community

ExecStart=/usr/bin/podman run --rm \
  --name certus-neo4j \
  --network certus \
  -p 7474:7474 \
  -p 7687:7687 \
  -v /var/lib/certus/data/neo4j:/data:z \
  -e NEO4J_AUTH=${NEO4J_AUTH} \
  -e NEO4J_server_memory_heap_initial__size=512m \
  -e NEO4J_server_memory_heap_max__size=2g \
  docker.io/library/neo4j:5-community

ExecStop=/usr/bin/podman stop certus-neo4j

[Install]
WantedBy=multi-user.target
