# deployment/shared/systemd/certus-redis.service
[Unit]
Description=Certus Redis Cache
After=network-online.target podman.service
Wants=network-online.target
Requires=podman.service

[Service]
Type=simple
User=certus
Restart=always
RestartSec=10s
TimeoutStartSec=60

EnvironmentFile=/opt/certus/.env

ExecStartPre=-/usr/bin/podman stop certus-redis
ExecStartPre=-/usr/bin/podman rm certus-redis
ExecStartPre=/usr/bin/podman pull docker.io/library/redis:7-alpine

ExecStart=/usr/bin/podman run --rm \
  --name certus-redis \
  --network certus \
  -p 6379:6379 \
  -v /var/lib/certus/data/redis:/data:z \
  docker.io/library/redis:7-alpine \
  redis-server --requirepass ${REDIS_PASSWORD} --appendonly yes

ExecStop=/usr/bin/podman stop certus-redis

[Install]
WantedBy=multi-user.target
